# komga2/rootfs/etc/services.d/komga/run
#!/usr/bin/with-contenv bashio

JAVAMEM=$(bashio::config 'maxmem')

# Run the komga server
cd /
exec java -Xmx100m -jar /komga.jar
