# karakeep_wg/rootfs/etc/services.d/meilisearch/run
#!/usr/bin/with-contenv bash
set -euo pipefail

# shellcheck source=/usr/lib/bashio/bashio.sh
source /usr/lib/bashio/bashio.sh

PUID=$(bashio::config 'PUID')
PGID=$(bashio::config 'PGID')

if [[ -f /data/karakeep.env ]]; then
  # shellcheck disable=SC1091
  source /data/karakeep.env
fi

if [[ -z "${MEILI_MASTER_KEY:-}" ]]; then
  bashio::log.error "MEILI_MASTER_KEY not set. Check /data/karakeep.env"
  bashio::exit.nok
fi

exec s6-setuidgid "${PUID}:${PGID}" meilisearch \
  --db-path /data/meilisearch \
  --http-addr 127.0.0.1:7700 \
  --master-key "${MEILI_MASTER_KEY}"
